# eudev/81-fuse.rules
#
# /lib/udev/rules.d/81-fuse.rules
# 

KERNEL=="fuse", ACTION=="add", RUN+="/bin/mount -t fusectl fusectl /sys/fs/fuse/connections"

# End of file.
